# auto-server.service
[Unit]
Description=myscript

[Service]
Type=simple
ExecStart=/home/dronspace/raspi-usbip-serv/server.sh

[Install]
WantedBy=multi-user.target
