# auto-ip.service
[Unit]
Description=myscript

[Service]
Type=simple
ExecStart=/home/dronspace/raspi-usbip-serv/autobind.sh

[Install]
WantedBy=multi-user.target
